# server/scripts/supervisor/mediamtx.conf
[program:mediamtx]
command = /home/username/mediamtx
directory = /home/username
user = username
startsecs = 0
autostart = true  
autorestart = true  
stdout_logfile = /home/username/logs/mediamtx.log  
stderr_logfile = /home/username/logs/mediamtx_err.log
